# etc/profile.d/floof-banner.sh
#!/bin/bash
echo "FloofOS [Fast Line-rate Offload On Fabric Operating System]"
echo "version: FloofOS 1.0"
