# usr/lib/systemd/system/vpp-auto-config.service
[Unit]
Description=VPP Auto Configuration
Before=vpp.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpp-auto-config.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
