# etc/systemd/system/nftables-dataplane.service
[Unit]
Description=nftables Firewall (Control Plane Protection)
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nft -f /etc/nftables.d/floofos.nft
ExecReload=/usr/sbin/nft -f /etc/nftables.d/floofos.nft
ExecStop=/usr/sbin/nft flush ruleset

[Install]
WantedBy=multi-user.target

